# Library/LaunchAgents/com.meta-manager.team-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meta-manager.team-sweep</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/riseof/environment/scripts/team-sweep.sh</string>
    </array>

    <!-- Run every 45 minutes -->
    <key>StartInterval</key>
    <integer>2700</integer>

    <key>StandardOutPath</key>
    <string>/Users/riseof/environment/.sweep.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/riseof/environment/.sweep.stderr.log</string>

    <!-- Run immediately after load, then on interval -->
    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
